# RSS-service/RSS-manager/supervisor/GetItems.conf
[program:laravel-queue-listen-item]
command=/usr/bin/php artisan queue:listen --queue=get_items
directory=/var/www/html
autostart=true
autorestart=true
numprocs=1
user=nginx
